# 1_compiled_reports/J70HPR/report.json
{"task_id": "20241230-J70HPR_1", "score": 7, "signatures": [{"name": "wrote_proc_memory", "short_description": "Wrote to the memory of another process", "description": "A process wrote data to an area of memory in another process.", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 12, "iocs": [{"truncated": false, "ioc": {"target_process": "cMd.exE", "target_process_id": 81, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"target_process": "cMd.exE", "target_process_id": 81, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"target_process": "cMd.exE", "target_process_id": 81, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 84, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 84, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 86, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 86, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 83, "process": "cMd.exE", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 83, "process": "cMd.exE", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 83, "process": "cMd.exE", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 86, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 84, "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}]}, "score": 1}, {"name": "susevent_adjustprivilegetoken", "short_description": "Suspicious behavior detected: AdjustPrivilegeToken", "description": "A privilege for an access token was adjusted. Windows uses the token to control the ability of the user to perform various system-related operations.", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 3, "iocs": [{"truncated": false, "ioc": {"process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}, {"truncated": false, "ioc": {"process": "shutdown.exe", "process_id": 84}}, {"truncated": false, "ioc": {"process": "shutdown.exe", "process_id": 84}}]}, "score": 1}, {"name": "files_drops_exe_general", "short_description": "Drops executable file", "description": "Drops executable file", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 1, "iocs": [{"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\Administrator\\1591121270.exe", "process": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "process_id": 80}}]}, "score": 6}, {"name": "registry_write_runkey", "short_description": "Add Run key to start application", "description": "Adds run entry to the registry to gain persistence across system reboots.", "ttps": ["T1547.001", "T1112"], "tags": [], "family": "", "iocs": {"truncated": false, "count": 2, "iocs": [{"truncated": false, "ioc": {"description": "Created registry key", "value": "\\REGISTRY\\USER\\S-1-5-21-4189341723-2754120001-2032404538-500\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run", "process": "reg.exe", "process_id": 83}}, {"truncated": false, "ioc": {"description": "Set registry key value", "value": "\\REGISTRY\\USER\\S-1-5-21-4189341723-2754120001-2032404538-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\1591121270", "process": "reg.exe", "process_id": 83}}]}, "score": 7}], "ttps": [{"id": "T1547.001", "name": "Registry Run Keys / Startup Folder", "tactics": ["Persistence", "Privilege Escalation"], "reference": "https://attack.mitre.org/techniques/T1547/001", "subtechniques": []}, {"id": "T1112", "name": "Modify Registry", "tactics": ["Defense Evasion"], "reference": "https://attack.mitre.org/techniques/T1112", "subtechniques": []}], "tags": [], "families": [], "processes": {"truncated": false, "count": 7, "process_list": [{"pid": 4876, "ppid": 3180, "procid": 80, "parent_procid": 59, "image": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "name": "VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe", "commandline": "\"C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_1b4c8e3b7482acef78331f81b96995c5.exe\"", "tracked": true, "injected": false, "state": "terminated", "start_ts": 578, "end_ts": 1485}, {"pid": 4908, "ppid": 4876, "procid": 81, "parent_procid": 80, "image": "C:\\Windows\\SysWOW64\\cMd.exE", "name": "cMd.exE", "commandline": "cMd.exE /c REG ADD HKEY_CURRENT_USER\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /v 1591121270 /t REG_SZ /d \"%userprofile%\\1591121270.exe\" /f", "tracked": true, "injected": false, "state": "terminated", "start_ts": 969, "end_ts": 1266}, {"pid": 4968, "ppid": 4908, "procid": 83, "parent_procid": 81, "image": "C:\\Windows\\SysWOW64\\reg.exe", "name": "reg.exe", "commandline": "REG  ADD HKEY_CURRENT_USER\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /v 1591121270 /t REG_SZ /d \"C:\\Users\\Administrator\\1591121270.exe\" /f", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1235, "end_ts": 1250}, {"pid": 5012, "ppid": 4876, "procid": 84, "parent_procid": 80, "image": "C:\\Windows\\SysWOW64\\shutdown.exe", "name": "shutdown.exe", "commandline": "\"C:\\Windows\\System32\\shutdown.exe\" /r /f /t 3", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1438, "end_ts": 1610}, {"pid": 5044, "ppid": 4876, "procid": 86, "parent_procid": 80, "image": "C:\\Windows\\SysWOW64\\cmd.exe", "name": "cmd.exe", "commandline": "\"C:\\Windows\\system32\\cmd.exe\" /c del C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VIRUSS~1.EXE > nul", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1453, "end_ts": 2531}, {"pid": 5072, "ppid": 512, "procid": 87, "parent_procid": 5, "image": "C:\\Windows\\system32\\wlrmdr.exe", "name": "wlrmdr.exe", "commandline": "-s -1 -f 2 -t You're about to be signed out -m Windows will shut down in less than a minute. -a 3", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1610, "end_ts": 5110}, {"pid": 4464, "ppid": 512, "procid": 95, "parent_procid": 5, "image": "C:\\Windows\\system32\\LogonUI.exe", "name": "LogonUI.exe", "commandline": "\"LogonUI.exe\" /flags:0x0 /state0:0xa3a20055 /state1:0x41c64e6d", "tracked": true, "injected": false, "state": "running", "start_ts": 4594, "end_ts": null}]}, "network": {"host": ["192.168.30.255", "8.8.8.8", "8.8.4.4", "224.0.0.252"], "udp": [{"ts": 1735530398.946021, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530398.946253, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530398.946365, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530399.727158, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530399.727487, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530399.727667, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530400.49498, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530401.259087, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530402.023356, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530402.780247, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530402.964036, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530403.180073, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735530403.545917, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735530403.578327, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735530403.661563, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735530403.717398, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530404.182428, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735530404.387321, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735530404.478833, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530404.666404, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735530405.057139, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735530405.184094, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735530405.387692, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735530405.536394, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530405.668768, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735530406.285628, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530406.395384, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735530407.05164, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530407.192999, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735530407.193144, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735530407.676695, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735530407.67689, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735530408.411274, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735530408.41148, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735530417.038103, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530417.038387, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530417.291395, "dstip": "224.0.0.252", "dstport": 5355, "srcip": "192.168.30.11", "srcport": 56678, "size": 26}, {"ts": 1735530417.419613, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530417.705358, "dstip": "224.0.0.252", "dstport": 5355, "srcip": "192.168.30.11", "srcport": 56678, "size": 26}, {"ts": 1735530417.781321, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530417.78151, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530418.155927, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530418.546833, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530418.546937, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530418.906112, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530419.312234, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530419.312627, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530419.67161, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530420.620844, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530420.732975, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735530421.368544, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530421.649273, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530422.11848, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530422.393911, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530422.893826, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735530423.159505, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530423.646727, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735530423.652051, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735530423.910081, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530424.690717, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530425.440859, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530426.683106, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530427.440724, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735530428.221996, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}], "dns": {"query": [{"ts": 1735530403.180073, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735530403.661563, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735530404.182428, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735530404.387321, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735530404.666404, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735530405.184094, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735530405.387692, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735530405.668768, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735530406.395384, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735530407.192999, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735530407.193144, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735530407.676695, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735530407.67689, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735530408.411274, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735530408.41148, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}]}}, "screenshot": [{"name": "9443.jpg", "percentage": 9.6675}, {"name": "11601.jpg", "percentage": 90.3325}]}