# 1_compiled_reports/RAW4SM/report.json
{"task_id": "20241231-RAW4SM_1", "score": 7, "signatures": [{"name": "wrote_proc_memory", "short_description": "Wrote to the memory of another process", "description": "A process wrote data to an area of memory in another process.", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 12, "iocs": [{"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 86, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 86, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 86, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cMd.exE", "target_process_id": 83, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 85, "process": "cMd.exE", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 85, "process": "cMd.exE", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "cMd.exE", "target_process_id": 83, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 85, "process": "cMd.exE", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 87, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 87, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 87, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cMd.exE", "target_process_id": 83, "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}]}, "score": 1}, {"name": "susevent_adjustprivilegetoken", "short_description": "Suspicious behavior detected: AdjustPrivilegeToken", "description": "A privilege for an access token was adjusted. Windows uses the token to control the ability of the user to perform various system-related operations.", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 3, "iocs": [{"truncated": false, "ioc": {"process": "shutdown.exe", "process_id": 86}}, {"truncated": false, "ioc": {"process": "shutdown.exe", "process_id": 86}}, {"truncated": false, "ioc": {"process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}]}, "score": 1}, {"name": "files_drops_exe_general", "short_description": "Drops executable file", "description": "Drops executable file", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 1, "iocs": [{"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\Administrator\\4083660.exe", "process": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "process_id": 82}}]}, "score": 6}, {"name": "registry_write_runkey", "short_description": "Add Run key to start application", "description": "Adds run entry to the registry to gain persistence across system reboots.", "ttps": ["T1547.001", "T1112"], "tags": [], "family": "", "iocs": {"truncated": false, "count": 2, "iocs": [{"truncated": false, "ioc": {"description": "Set registry key value", "value": "\\REGISTRY\\USER\\S-1-5-21-4189341723-2754120001-2032404538-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\4083660", "process": "reg.exe", "process_id": 85}}, {"truncated": false, "ioc": {"description": "Created registry key", "value": "\\REGISTRY\\USER\\S-1-5-21-4189341723-2754120001-2032404538-500\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run", "process": "reg.exe", "process_id": 85}}]}, "score": 7}], "ttps": [{"id": "T1547.001", "name": "Registry Run Keys / Startup Folder", "tactics": ["Persistence", "Privilege Escalation"], "reference": "https://attack.mitre.org/techniques/T1547/001", "subtechniques": []}, {"id": "T1112", "name": "Modify Registry", "tactics": ["Defense Evasion"], "reference": "https://attack.mitre.org/techniques/T1112", "subtechniques": []}], "tags": [], "families": [], "processes": {"truncated": false, "count": 7, "process_list": [{"pid": 4856, "ppid": 3180, "procid": 82, "parent_procid": 59, "image": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "name": "VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe", "commandline": "\"C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_7b5df649e489732b7df6fc37b4c8ee7c.exe\"", "tracked": true, "injected": false, "state": "terminated", "start_ts": 406, "end_ts": 1031}, {"pid": 4924, "ppid": 4856, "procid": 83, "parent_procid": 82, "image": "C:\\Windows\\SysWOW64\\cMd.exE", "name": "cMd.exE", "commandline": "cMd.exE /c REG ADD HKEY_CURRENT_USER\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /v 4083660 /t REG_SZ /d \"%userprofile%\\4083660.exe\" /f", "tracked": true, "injected": false, "state": "terminated", "start_ts": 640, "end_ts": 859}, {"pid": 4976, "ppid": 4924, "procid": 85, "parent_procid": 83, "image": "C:\\Windows\\SysWOW64\\reg.exe", "name": "reg.exe", "commandline": "REG  ADD HKEY_CURRENT_USER\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /v 4083660 /t REG_SZ /d \"C:\\Users\\Administrator\\4083660.exe\" /f", "tracked": true, "injected": false, "state": "terminated", "start_ts": 828, "end_ts": 843}, {"pid": 5036, "ppid": 4856, "procid": 86, "parent_procid": 82, "image": "C:\\Windows\\SysWOW64\\shutdown.exe", "name": "shutdown.exe", "commandline": "\"C:\\Windows\\System32\\shutdown.exe\" /r /f /t 3", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1000, "end_ts": 1203}, {"pid": 5048, "ppid": 4856, "procid": 87, "parent_procid": 82, "image": "C:\\Windows\\SysWOW64\\cmd.exe", "name": "cmd.exe", "commandline": "\"C:\\Windows\\system32\\cmd.exe\" /c del C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VIRUSS~1.EXE > nul", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1015, "end_ts": 1984}, {"pid": 5116, "ppid": 512, "procid": 90, "parent_procid": 5, "image": "C:\\Windows\\system32\\wlrmdr.exe", "name": "wlrmdr.exe", "commandline": "-s -1 -f 2 -t You're about to be signed out -m Windows will shut down in less than a minute. -a 3", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1203, "end_ts": 4500}, {"pid": 3292, "ppid": 512, "procid": 96, "parent_procid": 5, "image": "C:\\Windows\\system32\\LogonUI.exe", "name": "LogonUI.exe", "commandline": "\"LogonUI.exe\" /flags:0x0 /state0:0xa3a20055 /state1:0x41c64e6d", "tracked": true, "injected": false, "state": "running", "start_ts": 4218, "end_ts": null}]}, "network": {"host": ["192.168.30.255", "8.8.8.8", "8.8.4.4", "224.0.0.252"], "udp": [{"ts": 1735617728.947617, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617728.947813, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617728.947947, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617729.76015, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617729.760331, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617729.760412, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617730.517963, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617731.150606, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617731.277687, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617731.902748, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617731.997604, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735617732.034379, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617732.316319, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735617732.371786, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735617732.657926, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617732.798696, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617733.029745, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735617733.326357, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735617733.373317, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735617733.576762, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735617733.587499, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735617733.742741, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617734.029661, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735617734.338684, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735617734.385444, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735617734.4948, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617735.088097, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735617735.260026, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617735.774641, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735617736.042048, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735617736.042107, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735617736.354835, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735617736.35503, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 41}, {"ts": 1735617736.401486, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735617736.401629, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 42}, {"ts": 1735617736.588582, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735617738.103037, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735617745.938013, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617745.943251, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617746.482223, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617746.553974, "dstip": "224.0.0.252", "dstport": 5355, "srcip": "192.168.30.11", "srcport": 63203, "size": 26}, {"ts": 1735617746.698416, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617746.698524, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617746.965349, "dstip": "224.0.0.252", "dstport": 5355, "srcip": "192.168.30.11", "srcport": 63203, "size": 26}, {"ts": 1735617747.245403, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617747.448606, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617747.448814, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617747.995435, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617748.198438, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617748.198593, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617748.745357, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617749.525518, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617749.787131, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735617750.273608, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617751.023708, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617751.49717, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617751.778243, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735617752.24703, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617752.555304, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735617752.559027, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735617753.012572, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617753.794013, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617754.043916, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735617754.543905, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617755.293852, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617755.55949, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735617756.418341, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617757.07512, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735617757.153186, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617757.903209, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735617760.778362, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 191}, {"ts": 1735617761.778324, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 191}, {"ts": 1735617774.908698, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}], "dns": {"query": [{"ts": 1735617731.997604, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735617732.316319, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735617732.371786, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735617733.029745, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735617733.326357, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735617733.373317, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735617734.029661, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735617734.338684, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735617734.385444, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735617736.042048, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735617736.042107, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735617736.354835, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735617736.35503, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735617736.401486, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735617736.401629, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "sls.update.microsoft.com"}]}}, "screenshot": [{"name": "9245.jpg", "percentage": 8.595}, {"name": "10314.jpg", "percentage": 0.8825}, {"name": "11373.jpg", "percentage": 90.5225}]}