# 1_compiled_reports/TJKF98/report.json
{"task_id": "20241227-TJKF98_1", "score": 7, "signatures": [{"name": "wrote_proc_memory", "short_description": "Wrote to the memory of another process", "description": "A process wrote data to an area of memory in another process.", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 12, "iocs": [{"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 83, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 83, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 87, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 87, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 86, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 85, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 85, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 85, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 83, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 86, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 86, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 87, "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}]}, "score": 1}, {"name": "susevent_adjustprivilegetoken", "short_description": "Suspicious behavior detected: AdjustPrivilegeToken", "description": "A privilege for an access token was adjusted. Windows uses the token to control the ability of the user to perform various system-related operations.", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 3, "iocs": [{"truncated": false, "ioc": {"process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}, {"truncated": false, "ioc": {"process": "shutdown.exe", "process_id": 86}}, {"truncated": false, "ioc": {"process": "shutdown.exe", "process_id": 86}}]}, "score": 1}, {"name": "files_drops_exe_general", "short_description": "Drops executable file", "description": "Drops executable file", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 1, "iocs": [{"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\Administrator\\26419606.exe", "process": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "process_id": 82}}]}, "score": 6}, {"name": "registry_write_runkey", "short_description": "Add Run key to start application", "description": "Adds run entry to the registry to gain persistence across system reboots.", "ttps": ["T1547.001", "T1112"], "tags": [], "family": "", "iocs": {"truncated": false, "count": 2, "iocs": [{"truncated": false, "ioc": {"description": "Set registry key value", "value": "\\REGISTRY\\USER\\S-1-5-21-4189341723-2754120001-2032404538-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\26419606", "process": "reg.exe", "process_id": 85}}, {"truncated": false, "ioc": {"description": "Created registry key", "value": "\\REGISTRY\\USER\\S-1-5-21-4189341723-2754120001-2032404538-500\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run", "process": "reg.exe", "process_id": 85}}]}, "score": 7}], "ttps": [{"id": "T1547.001", "name": "Registry Run Keys / Startup Folder", "tactics": ["Persistence", "Privilege Escalation"], "reference": "https://attack.mitre.org/techniques/T1547/001", "subtechniques": []}, {"id": "T1112", "name": "Modify Registry", "tactics": ["Defense Evasion"], "reference": "https://attack.mitre.org/techniques/T1112", "subtechniques": []}], "tags": [], "families": [], "processes": {"truncated": false, "count": 7, "process_list": [{"pid": 4792, "ppid": 3180, "procid": 82, "parent_procid": 59, "image": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "name": "VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe", "commandline": "\"C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_2fdf9dec8c297f24516fdb4fa6d30898.exe\"", "tracked": true, "injected": false, "state": "terminated", "start_ts": 547, "end_ts": 1390}, {"pid": 4936, "ppid": 4792, "procid": 83, "parent_procid": 82, "image": "C:\\Windows\\SysWOW64\\cmd.exe", "name": "cmd.exe", "commandline": "cmd.exe /c REG ADD HKEY_CURRENT_USER\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /v 26419606 /t REG_SZ /d \"%userprofile%\\26419606.exe\" /f", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1078, "end_ts": 1265}, {"pid": 4984, "ppid": 4936, "procid": 85, "parent_procid": 83, "image": "C:\\Windows\\SysWOW64\\reg.exe", "name": "reg.exe", "commandline": "REG  ADD HKEY_CURRENT_USER\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /v 26419606 /t REG_SZ /d \"C:\\Users\\Administrator\\26419606.exe\" /f", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1172, "end_ts": 1265}, {"pid": 5032, "ppid": 4792, "procid": 86, "parent_procid": 82, "image": "C:\\Windows\\SysWOW64\\shutdown.exe", "name": "shutdown.exe", "commandline": "\"C:\\Windows\\System32\\shutdown.exe\" /r /f /t 3", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1359, "end_ts": 1562}, {"pid": 5044, "ppid": 4792, "procid": 87, "parent_procid": 82, "image": "C:\\Windows\\SysWOW64\\cmd.exe", "name": "cmd.exe", "commandline": "\"C:\\Windows\\system32\\cmd.exe\" /c del C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VIRUSS~1.EXE > nul", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1375, "end_ts": 2312}, {"pid": 5108, "ppid": 512, "procid": 90, "parent_procid": 5, "image": "C:\\Windows\\system32\\wlrmdr.exe", "name": "wlrmdr.exe", "commandline": "-s -1 -f 2 -t You're about to be signed out -m Windows will shut down in less than a minute. -a 3", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1562, "end_ts": 4859}, {"pid": 3292, "ppid": 512, "procid": 96, "parent_procid": 5, "image": "C:\\Windows\\system32\\LogonUI.exe", "name": "LogonUI.exe", "commandline": "\"LogonUI.exe\" /flags:0x0 /state0:0xa3a20055 /state1:0x41c64e6d", "tracked": true, "injected": false, "state": "running", "start_ts": 4547, "end_ts": null}]}, "network": {"host": ["192.168.30.255", "8.8.8.8", "8.8.4.4", "224.0.0.252"], "udp": [{"ts": 1735338297.141642, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338297.141836, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338297.141953, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338297.892018, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338297.892198, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338297.892319, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338298.648262, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338299.403625, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338299.506967, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338299.928151, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735338300.159021, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338300.252694, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338300.334917, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 42}, {"ts": 1735338300.372869, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 41}, {"ts": 1735338300.923378, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338300.954997, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735338301.007351, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338301.335798, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 42}, {"ts": 1735338301.367086, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 41}, {"ts": 1735338301.680472, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735338301.684851, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735338301.964359, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735338302.086163, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338302.343956, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 42}, {"ts": 1735338302.375552, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 41}, {"ts": 1735338302.828337, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338303.181339, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735338303.587036, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338303.835201, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735338303.97899, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735338303.979055, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735338304.353409, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 42}, {"ts": 1735338304.353633, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "size": 42}, {"ts": 1735338304.384456, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 41}, {"ts": 1735338304.384607, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "size": 41}, {"ts": 1735338304.680834, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735338306.185994, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735338313.961347, "dstip": "224.0.0.252", "dstport": 5355, "srcip": "192.168.30.11", "srcport": 52490, "size": 26}, {"ts": 1735338314.030072, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338314.030409, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338314.030698, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338314.378964, "dstip": "224.0.0.252", "dstport": 5355, "srcip": "192.168.30.11", "srcport": 52490, "size": 26}, {"ts": 1735338314.779478, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338314.779583, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338314.779636, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338315.545197, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338315.545337, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338315.545391, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338316.310845, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338316.310979, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338316.311063, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338317.087367, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338317.827846, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338318.609283, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338318.91709, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338319.365356, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735338319.677739, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338320.116704, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735338320.120821, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735338320.427815, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338321.17827, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338321.63098, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735338321.943425, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338322.709068, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338323.131049, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735338323.864391, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338324.63196, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338324.646703, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735338325.380899, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735338327.771581, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 191}, {"ts": 1735338328.787243, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 191}, {"ts": 1735338344.259345, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735338408.137856, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}], "dns": {"query": [{"ts": 1735338299.928151, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735338300.334917, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735338300.372869, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735338300.954997, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735338301.335798, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735338301.367086, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735338301.964359, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735338302.343956, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735338302.375552, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735338303.97899, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735338303.979055, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735338304.353409, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735338304.353633, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 57059, "type": "A", "name": "sls.update.microsoft.com"}, {"ts": 1735338304.384456, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "insiderppe.cloudapp.net"}, {"ts": 1735338304.384607, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 63674, "type": "A", "name": "insiderppe.cloudapp.net"}]}}, "screenshot": [{"name": "8348.jpg", "percentage": 8.710833333333333}, {"name": "10453.jpg", "percentage": 91.28916666666667}]}