# 1_compiled_reports/PKPW5C/report.json
{"task_id": "20241230-PKPW5C_1", "score": 8, "signatures": [{"name": "wrote_proc_memory", "short_description": "Wrote to the memory of another process", "description": "A process wrote data to an area of memory in another process.", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 21, "iocs": [{"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 83, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 83, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 83, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 86, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 86, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 86, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "taskkill.exe", "target_process_id": 89, "process": "cmd.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "taskkill.exe", "target_process_id": 89, "process": "cmd.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "taskkill.exe", "target_process_id": 89, "process": "cmd.exe", "process_id": 82}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 82, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 82, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "cmd.exe", "target_process_id": 82, "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 90, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 90, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 90, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 91, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 91, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "reg.exe", "target_process_id": 91, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 92, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 92, "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"target_process": "shutdown.exe", "target_process_id": 92, "process": "cmd.exe", "process_id": 83}}]}, "score": 1}, {"name": "susevent_adjustprivilegetoken", "short_description": "Suspicious behavior detected: AdjustPrivilegeToken", "description": "A privilege for an access token was adjusted. Windows uses the token to control the ability of the user to perform various system-related operations.", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 3, "iocs": [{"truncated": false, "ioc": {"process": "shutdown.exe", "process_id": 92}}, {"truncated": false, "ioc": {"process": "shutdown.exe", "process_id": 92}}, {"truncated": false, "ioc": {"process": "taskkill.exe", "process_id": 89}}]}, "score": 1}, {"name": "files_drops_exe_user", "short_description": "Drops file in user directory", "description": "Drops file in user directory", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 6, "iocs": [{"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\CFSetup325.exe", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\000.bat", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\000.bat", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\System38.bat", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\CFSetup325.exe", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\System38.bat", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}]}, "score": 6}, {"name": "files_drops_exe_general", "short_description": "Drops executable file", "description": "Drops executable file", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 9, "iocs": [{"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Windows\\System38.bat", "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Windows\\System38.bat", "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Windows\\System38.bat", "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\CFSetup325.exe", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\000.bat", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\000.bat", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\System38.bat", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\CFSetup325.exe", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}, {"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\System38.bat", "process": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "process_id": 81}}]}, "score": 6}, {"name": "cmd_taskkill_process", "short_description": "Kills process with taskkill.exe", "description": "Kills process with taskkill.exe", "ttps": ["T1562.001"], "tags": [], "family": "", "iocs": {"truncated": false, "count": 1, "iocs": [{"truncated": false, "ioc": {"description": "Created process", "value": "taskkill  /im explorer.exe /f ", "process": "taskkill.exe", "process_id": 89}}]}, "score": 7}, {"name": "registry_disables_taskmgr", "short_description": "Disables Windows Task Manager", "description": "Modifies registry key(s) to prevent use of Task Manager.", "ttps": ["T1112", "T1562.001"], "tags": [], "family": "", "iocs": {"truncated": false, "count": 1, "iocs": [{"truncated": false, "ioc": {"description": "Set registry key value", "value": "\\REGISTRY\\USER\\S-1-5-21-4189341723-2754120001-2032404538-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\DisableTaskMgr", "process": "reg.exe", "process_id": 90}}]}, "score": 8}, {"name": "registry_write_runkey", "short_description": "Add Run key to start application", "description": "Adds run entry to the registry to gain persistence across system reboots.", "ttps": ["T1547.001", "T1112"], "tags": [], "family": "", "iocs": {"truncated": false, "count": 2, "iocs": [{"truncated": false, "ioc": {"description": "Created registry key", "value": "\\REGISTRY\\MACHINE\\Software\\WOW6432Node\\Microsoft\\Windows\\CurrentVersion\\Run", "process": "reg.exe", "process_id": 91}}, {"truncated": false, "ioc": {"description": "Set registry key value", "value": "\\REGISTRY\\MACHINE\\SOFTWARE\\WOW6432Node\\Microsoft\\Windows\\CurrentVersion\\Run\\Sys38", "process": "reg.exe", "process_id": 91}}]}, "score": 7}, {"name": "file_drops_windows", "short_description": "Drops file in Windows directory", "description": "Drops file in Windows directory", "ttps": [], "tags": [], "family": "", "iocs": {"truncated": false, "count": 3, "iocs": [{"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Windows\\System38.bat", "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"description": "Created file", "value": "C:\\Windows\\System38.bat", "process": "cmd.exe", "process_id": 83}}, {"truncated": false, "ioc": {"description": "Opened file for modifying", "value": "C:\\Windows\\System38.bat", "process": "cmd.exe", "process_id": 83}}]}, "score": 4}], "ttps": [{"id": "T1547.001", "name": "Registry Run Keys / Startup Folder", "tactics": ["Persistence", "Privilege Escalation"], "reference": "https://attack.mitre.org/techniques/T1547/001", "subtechniques": []}, {"id": "T1562.001", "name": "Disable or Modify Tools", "tactics": ["Defense Evasion"], "reference": "https://attack.mitre.org/techniques/T1562/001", "subtechniques": []}, {"id": "T1112", "name": "Modify Registry", "tactics": ["Defense Evasion"], "reference": "https://attack.mitre.org/techniques/T1112", "subtechniques": []}], "tags": [], "families": [], "processes": {"truncated": false, "count": 9, "process_list": [{"pid": 4884, "ppid": 3180, "procid": 81, "parent_procid": 59, "image": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "name": "VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe", "commandline": "\"C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe\"", "tracked": true, "injected": false, "state": "terminated", "start_ts": 407, "end_ts": 938}, {"pid": 4944, "ppid": 4884, "procid": 82, "parent_procid": 81, "image": "C:\\Windows\\SysWOW64\\cmd.exe", "name": "cmd.exe", "commandline": "C:\\Windows\\system32\\cmd.exe /c \"\"C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\System38.bat\" \"", "tracked": true, "injected": false, "state": "terminated", "start_ts": 844, "end_ts": 1922}, {"pid": 4952, "ppid": 4884, "procid": 83, "parent_procid": 81, "image": "C:\\Windows\\SysWOW64\\cmd.exe", "name": "cmd.exe", "commandline": "C:\\Windows\\system32\\cmd.exe /c \"\"C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\000.bat\" \"", "tracked": true, "injected": false, "state": "terminated", "start_ts": 860, "end_ts": 1657}, {"pid": 4984, "ppid": 4884, "procid": 86, "parent_procid": 81, "image": "C:\\Windows\\SysWOW64\\cmd.exe", "name": "cmd.exe", "commandline": "\"C:\\Windows\\system32\\cmd.exe\" /c del \"C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\VirusShare_5ce38e091ed9b4bb7052190e0159d000.exe\" >> NUL", "tracked": true, "injected": false, "state": "terminated", "start_ts": 891, "end_ts": 1204}, {"pid": 5092, "ppid": 4944, "procid": 89, "parent_procid": 82, "image": "C:\\Windows\\SysWOW64\\taskkill.exe", "name": "taskkill.exe", "commandline": "taskkill  /im explorer.exe /f ", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1219, "end_ts": 1610}, {"pid": 5108, "ppid": 4952, "procid": 90, "parent_procid": 83, "image": "C:\\Windows\\SysWOW64\\reg.exe", "name": "reg.exe", "commandline": "reg  add \"HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\" /v DisableTaskMgr /t REG_DWORD /d 1 /f", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1250, "end_ts": 1297}, {"pid": 3280, "ppid": 4952, "procid": 91, "parent_procid": 83, "image": "C:\\Windows\\SysWOW64\\reg.exe", "name": "reg.exe", "commandline": "reg  add \"HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\" /v Sys38 /t REG_SZ /d C:\\Windows\\System38.bat /f", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1313, "end_ts": 1344}, {"pid": 1004, "ppid": 4952, "procid": 92, "parent_procid": 83, "image": "C:\\Windows\\SysWOW64\\shutdown.exe", "name": "shutdown.exe", "commandline": "Shutdown  /t 0 /s", "tracked": true, "injected": false, "state": "terminated", "start_ts": 1407, "end_ts": 1469}, {"pid": 2908, "ppid": 512, "procid": 94, "parent_procid": 5, "image": "C:\\Windows\\system32\\LogonUI.exe", "name": "LogonUI.exe", "commandline": "\"LogonUI.exe\" /flags:0x0 /state0:0xa3a1b855 /state1:0x41c64e6d", "tracked": true, "injected": false, "state": "running", "start_ts": 1454, "end_ts": null}]}, "network": {"host": ["192.168.30.255", "8.8.8.8", "8.8.4.4"], "udp": [{"ts": 1735534514.281266, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534514.281464, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534514.281605, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534515.062415, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534515.062578, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534515.062658, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534515.843831, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534516.599913, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534517.355405, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534518.110322, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 68}, {"ts": 1735534518.299117, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735534518.861789, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735534518.877497, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735534519.048974, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735534519.166725, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735534519.799302, "dstip": "192.168.30.255", "dstport": 137, "srcip": "192.168.30.11", "srcport": 137, "size": 50}, {"ts": 1735534520.160015, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735534520.228508, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 201}, {"ts": 1735534520.362112, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735534521.174519, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735534521.861535, "dstip": "192.168.30.255", "dstport": 138, "srcip": "192.168.30.11", "srcport": 138, "size": 179}, {"ts": 1735534523.183202, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}, {"ts": 1735534523.183366, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "size": 49}], "dns": {"query": [{"ts": 1735534519.166725, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735534520.160015, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735534521.174519, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735534523.183202, "dstip": "8.8.8.8", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}, {"ts": 1735534523.183366, "dstip": "8.8.4.4", "dstport": 53, "srcip": "192.168.30.11", "srcport": 62434, "type": "A", "name": "settings-win.data.microsoft.com"}]}}, "screenshot": [{"name": "10298.jpg", "percentage": 91.41833333333334}]}